# OnlineJudge/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:cache]
command=redis-server
user=root

[program:judge]
command=java -jar online-judge.jar
user=root
